# molly/supervisor.cfg
[unix_http_server]
file=/var/lib/molly/supervisord.sock
chown=molly:molly

[supervisord]
logfile=/var/log/molly/supervisord.log
pidfile=/var/lib/molly/supervisord.pid
user=molly

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/lib/molly/supervisor.sock

[program:mollyrest]
directory=/opt/molly/
command=/opt/molly/bin/mollyrest start
autostart=true
autorestart=true
redirect_stderr=true
user=molly

[program:mollyui]
directory=/opt/molly/
command=/opt/molly/bin/mollyui start
autostart=true
autorestart=true
redirect_stderr=true
user=molly