# molly/supervisor-debug.cfg
[unix_http_server]
file=/var/lib/molly/supervisord.sock
chown=molly:molly

[supervisord]
logfile=/var/log/molly/supervisord.log
pidfile=/var/lib/molly/supervisord.pid
user=molly

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/lib/molly/supervisord.sock

[program:mollyrest]
directory=/opt/molly/
command=/opt/molly/bin/mollyrest start_debug -a 0.0.0.0
autostart=true
autorestart=true
redirect_stderr=true
user=molly

[program:mollyui]
directory=/opt/molly/
command=/opt/molly/bin/mollyui start_debug -a 0.0.0.0
autostart=true
autorestart=true
redirect_stderr=true
user=molly
environment=MOLLY_UI_SETTINGS=/vagrant/conf/ui.py

[program:mollytaskbeat]
directory=/opt/molly/
command=/opt/molly/bin/mollyrest taskbeat
autostart=true
autorestart=true
redirect_stderr=true
user=molly

[program:mollytaskworker]
directory=/opt/molly/
command=/opt/molly/bin/mollyrest taskworker
autostart=true
autorestart=true
redirect_stderr=true
user=molly